# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/var/run/supervisord.pid

[program:api]
command=node /app/api/dist/index.js
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api.log
stdout_logfile_maxbytes=10MB
environment=NODE_ENV="production"

[program:worker]
command=node /app/api/dist/workers/index.js
directory=/app
autostart=true
autorestart=true
startsecs=5
startretries=5
stdout_logfile=/var/log/supervisor/worker.log
stderr_logfile=/var/log/supervisor/worker.log
stdout_logfile_maxbytes=10MB
environment=NODE_ENV="production"

[program:web]
command=node /app/web/server.js
directory=/app/web
autostart=true
autorestart=true
startsecs=10
startretries=5
stdout_logfile=/var/log/supervisor/web.log
stderr_logfile=/var/log/supervisor/web.log
stdout_logfile_maxbytes=10MB
environment=NODE_ENV="production",PORT="3000",HOSTNAME="0.0.0.0",INTERNAL_API_URL="http://127.0.0.1:4000"

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
stdout_logfile_maxbytes=5MB
